# stacks/terminal/files/devvm/ttyd.service
[Unit]
Description=ttyd Terminal Service (multi-session lobby + per-Authentik-user attach on port 7681)
After=network.target

[Service]
ExecStart=/usr/local/bin/ttyd -W -a -H X-authentik-username -t enableClipboard=true -I /usr/local/share/ttyd/index.html -p 7681 /usr/local/bin/tmux-attach.sh
Restart=always
User=wizard

[Install]
WantedBy=multi-user.target
