# scripts/daily-backup.service
[Unit]
Description=Daily backup: PVC snapshots + SQLite + pfsense to sda
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/daily-backup
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-backup
TimeoutStartSec=3600
