# scripts/lvm-pvc-snapshot.timer
[Unit]
Description=Daily LVM thin snapshots of Proxmox CSI PVCs

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
