# scripts/nfs-change-tracker.service
[Unit]
Description=Track NFS filesystem changes for incremental offsite backup
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/inotifywait -m -r \
    --format '%%w%%f' \
    -e create -e modify -e moved_to -e delete \
    --exclude '(/\..*swp$|/\.nfs|/\.Trash|\.db-shm$|\.db-wal$|\.db-journal$|/stats/.*\.stat$)' \
    /srv/nfs \
    /srv/nfs-ssd
StandardOutput=append:/mnt/backup/.nfs-changes.log
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
