# scripts/nfs-mirror.timer
[Unit]
Description=Weekly local NFS mirror to /mnt/backup

[Timer]
OnCalendar=Mon *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
