# scripts/offsite-sync-backup.timer
[Unit]
Description=Daily offsite sync: sda + NFS changes to Synology

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
