# scripts/t3-autoupdate.service
[Unit]
Description=Track latest t3 nightly (health-checked, idle-only restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/t3-autoupdate
