# scripts/t3-autoupdate.timer
[Unit]
Description=Daily gated t3 nightly tracker (health-checked + canary + auto-rollback)

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=1h
# Persistent deliberately OMITTED: this now installs a NEW build + migrates DBs +
# restarts serves, so a missed 04:00 run must NOT fire on boot mid-day with users
# active (a 2026-06-09 contributing factor). Skipping a day is fine — the next
# 04:00 picks up the latest nightly.

[Install]
WantedBy=timers.target
