# scripts/t3-provision-users.timer
[Unit]
Description=Periodic t3 per-user reconcile

[Timer]
OnBootSec=2min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
