# scripts/vault-token-renew.timer
[Unit]
Description=Daily renewal of the periodic Vault token in ~/.vault-token

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
