# scripts/weekly-backup.service
[Unit]
Description=Daily backup: PVC snapshots + SQLite + pfsense to sda
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/weekly-backup
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weekly-backup
TimeoutStartSec=3600
