# scripts/weekly-backup.timer
[Unit]
Description=Daily backup: PVC snapshots + SQLite + pfsense to sda

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
